# etc/supervisor.d/mysqld.ini
[program:mysqld]
command=/usr/bin/mysqld
redirect_stderr=true
user=mysql
