# etc/supervisor.d/mysqld.ini
[program:mysqld]
command=/usr/bin/mysqld --user=mysql --datadir=/var/lib/mysql
redirect_stderr=true
