# etc/supervisor.d/redis-server.ini
[program:redis-server]
command=/usr/bin/redis-server
redirect_stderr=true
