# rootfs/etc/supervisor.d/mysqld.ini
[program:mysqld]
command=/usr/bin/pidproxy /var/run/mysqld/mysqld.pid /usr/bin/mysqld_safe --pid-file=/var/run/mysqld/mysqld.pid
redirect_stderr=true
